# iot/supervisor/coap.conf
[program:coap]
command=coap-server -A 0.0.0.0
autostart=true
autorestart=true
stderr_logfile=/var/log/coap.err.log
stdout_logfile=/var/log/coap.out.log