# iot/supervisor/mosquitto.conf
[program:mosquitto]
command=/usr/sbin/mosquitto -c /var/conf/mosquitto.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/mosquitto.err.log
stdout_logfile=/var/log/mosquitto.out.log