# containers/iot/supervisor/cupsd.conf
[program:cupsd]
command=/usr/sbin/cupsd -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cups.err.log
stdout_logfile=/var/log/cups.out.log