# containers/iot/supervisor/inetd.conf
[program:inetd]
command=/usr/sbin/inetd -f /etc/inetd/inetd.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/inetd.err.log
stdout_logfile=/var/log/inetd.out.log