# containers/iot/supervisor/post.conf
[program:post]
command = /post.sh
startsecs = 0
autorestart = false
startretries = 1
priority=1000