# containers/iot/supervisor/vsftpd.conf
[program:vsftpd]
command=/usr/sbin/vsftpd
autostart=true
autorestart=true
stderr_logfile=/var/log/ftp.err.log
stdout_logfile=/var/log/ftp.out.log